# datamart/files/supervisor.conf
[supervisord]
nodaemon=true

[program:datamart]
command=bash /datamart | tail -f 
