# datamart/files/supervisor.conf
[supervisord]
nodaemon=true

[program:datamart]
command=/usr/bin/tail -f /var/log/messages
